# resources/systemd/devidservice.service
[Unit]
Description=Provisioning Wi-Fi Support Service
After=bb-wl18xx-wlan0.service
Before=webapp.service

[Service]
Restart=on-failure
ExecStart=/bin/sh /opt/apispoint/provision/services/devidservice.sh

[Install]
WantedBy=multi-user.target